# crontab
15 7 * * * python3 /config/vaccine_script.py
